# connectnfs.service
[Unit]
Description=Connect NFS
After=network.target logitechmediaserver-git.service

[Service]
ExecStart=/opt/getip.sh
Restart=always

[Install]
WantedBy=multi-user.target
