# connectnfs.service
[Unit]
Description=Connect NFS
After=network.target logitechmediaserver-git.service

[Service]
Type=simple
ExecStart=/opt/getip.sh

[Install]
WantedBy=multi-user.target
