# connectnfs.service
[Unit]
Description=Connect NFS
After=network.target logitechmediaserver-git.service

Type=simple
ExecStartPre=/bin/sleep 10
ExecStart=/opt/getip.sh

[Install]
WantedBy=multi-user.target
